# roles/backup/files/backup.timer
[Unit]
Description=Run backup.service every month

[Timer]
OnCalendar=*-*-01 02:00:00

[Install]
WantedBy=timers.target
