# roles/backup/files/backup.service
[Unit]
Description=Backup all data

[Service]
Type=simple
ExecStart=/opt/backup.sh
