# roles/nextcloud/files/crontab
10 04 10 01 * cd /opt/nextcloud; sudo -u nobody /usr/bin/php7.3 occ trashbin:cleanup rickou_412