# roles/nextcloud/files/nextcloudTrash.service
[Unit]
Description=Nextcloud Trash

[Service]
Type=simple
ExecStart=/opt/nextcloudTrash.sh
