# roles/nextcloud/files/nextcloudCronAppAutoUpdate.service
[Unit]
Description=Nextcloud Trash

[Service]
Type=simple
ExecStart=/opt/nextcloudAppAutoUpdate.sh
