# roles/checkVersion/files/checkVersion.service
[Unit]
Description=checkVersion

[Service]
Type=simple
ExecStart=/opt/checkVersion.sh
User=ovhuser


