# roles/checkVersion/files/checkVersion.timer
[Unit]
Description=Run checkVersion.service every day

[Timer]
OnCalendar=Mon *-*-* 05:00:00

[Install]
WantedBy=timers.target
